# Documents/GPIO_fix/ICA_99-gpio-sysfs-perms.rules
ACTION=="add|change", SUBSYSTEM=="gpio", KERNEL=="gpio[0-9]*", RUN+="/usr/local/sbin/ICA_fix-gpio-perms.sh %k"
